# engine/scheduler/com.mikai.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mikai.daily-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/briancho/Desktop/MIKAI/engine/scheduler/daily-sync.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/briancho/Desktop/MIKAI/engine/scheduler/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/briancho/Desktop/MIKAI/engine/scheduler/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/briancho/Desktop/MIKAI</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
